# com.bixbit.mining-guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Identity -->
    <key>Label</key>
    <string>com.bixbit.mining-guardian</string>

    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/BigBobby/Documents/GitHub/Mining Gaurdian/start_guardian.sh</string>
    </array>

    <!-- Working directory (so config.json and guardian.db are found) -->
    <key>WorkingDirectory</key>
    <string>/Users/BigBobby/Documents/GitHub/Mining Gaurdian</string>

    <!-- Environment variables (credentials loaded from .env manually below) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>GUARDIAN_CONFIG</key>
        <string>/Users/BigBobby/Documents/GitHub/Mining Gaurdian/config.json</string>
    </dict>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle restarts — wait 30s before restarting after a crash -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Log stdout and stderr to files -->
    <key>StandardOutPath</key>
    <string>/Users/BigBobby/Documents/GitHub/Mining Gaurdian/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/BigBobby/Documents/GitHub/Mining Gaurdian/logs/launchd_stderr.log</string>

</dict>
</plist>
