# deploy/intelligence-report.service
[Unit]
Description=Mining Guardian Intelligence Report API
After=network.target dashboard-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/Mining-Guardian
ExecStart=/root/Mining-Guardian/venv/bin/python /root/Mining-Guardian/api/intelligence_report_api.py
Environment=GUARDIAN_DB=/root/Mining-Guardian/guardian.db
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
