# deploy/slack-commands.service
[Unit]
Description=Mining Guardian Slack Command Handler
After=network.target mining-guardian.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/Mining-Guardian
ExecStart=/root/Mining-Guardian/venv/bin/python api/slack_command_handler.py
Restart=always
RestartSec=10
EnvironmentFile=/root/Mining-Guardian/.env

[Install]
WantedBy=multi-user.target
