# deploy/slack-listener.service
[Unit]
Description=Mining Guardian Slack Approval Listener
After=network.target mining-guardian.service

[Service]
EnvironmentFile=/root/Mining-Gaurdian/.env
Type=simple
User=root
WorkingDirectory=/root/Mining-Gaurdian
ExecStart=/root/Mining-Gaurdian/venv/bin/python api/slack_approval_listener.py
Restart=always
RestartSec=10
Environment=PATH=/root/Mining-Gaurdian/venv/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=multi-user.target
