# installer/macos-pkg/resources/launchd/com.miningguardian.slack-commands.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.miningguardian.slack-commands.plist (Bucket 6 / installer rebuild)

  launchd LaunchDaemon for the Slack slash-command handler on macOS.
  Mirrors the systemd unit `slack-commands.service` (deploy/slack-commands.service)
  running today on the VPS.

  Service entry point: api/slack_command_handler.py
  Restart policy:      always restart on crash, 10s throttle.

  Install path: copied by postinstall.sh to
      /Library/LaunchDaemons/com.miningguardian.slack-commands.plist
  and loaded with:
      sudo launchctl bootstrap system <path-to-plist>

  Path + env conventions (match scanner / dashboard-api / approval-api / feedback-loop):
    install root              /Library/Application Support/MiningGuardian
    .env (private secrets)    /Library/Application Support/MiningGuardian/.env  (mode 0600)
    venv                      /Library/Application Support/MiningGuardian/venv
    launcher wrapper          /Library/Application Support/MiningGuardian/bin/slack_commands_launcher.sh
    logs                      /Library/Application Support/MiningGuardian/logs/
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miningguardian.slack-commands</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/Application Support/MiningGuardian/bin/slack_commands_launcher.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Library/Application Support/MiningGuardian</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/slack-commands.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/slack-commands.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Library/Application Support/MiningGuardian/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
