# installer/macos-pkg/resources/launchd/com.miningguardian.dashboard-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.miningguardian.dashboard-api.plist (PR I-2 / Bucket-3 installer)

  launchd LaunchDaemon for the Mining Guardian dashboard API (the FastAPI
  service that fronts the Grafana embed + the operator dashboard pages).
  Mirrors the systemd unit `dashboard-api.service` on ROBS-PC.

  Install path: copied by postinstall.sh to
      /Library/LaunchDaemons/com.miningguardian.dashboard-api.plist

  Same path + env conventions as scanner.plist + feedback-loop-daemon.plist
  (see PR #41 for the full discussion of the .env-sourcing wrapper).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miningguardian.dashboard-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/Application Support/MiningGuardian/bin/dashboard_api_launcher.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Library/Application Support/MiningGuardian</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Standard</string>

    <key>StandardOutPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/dashboard_api.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/dashboard_api.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Library/Application Support/MiningGuardian/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
