# installer/macos-pkg/resources/launchd/com.miningguardian.intelligence-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.miningguardian.intelligence-report.plist (Bucket 6 / installer rebuild)

  launchd LaunchDaemon for the Intelligence Report API on macOS.
  Mirrors the systemd unit `intelligence-report.service`
  (deploy/intelligence-report.service) running today on the VPS.

  Service entry point: api/intelligence_report_api.py
  Restart policy:      always restart on crash, 5s throttle (matches systemd RestartSec=5).

  Bind/port hardening (S-8): the API binds to 127.0.0.1 by default, with
  CORS controlled via INTELLIGENCE_REPORT_CORS_ALLOW in .env (PR #66).
  This plist does not need to do anything special — the application reads
  .env and does the right thing.

  GUARDIAN_DB note: the systemd unit sets GUARDIAN_DB=/root/Mining-Guardian/guardian.db,
  but on macOS the install root is /usr/local/MiningGuardian, so we set
  GUARDIAN_DB=/usr/local/MiningGuardian/guardian.db here. The launcher wrapper
  also sources .env, which can override this if the operator wants a custom
  path.

  Install path: copied by postinstall.sh to
      /Library/LaunchDaemons/com.miningguardian.intelligence-report.plist
  and loaded with:
      sudo launchctl bootstrap system <path-to-plist>

  Path + env conventions (match scanner / dashboard-api / approval-api / feedback-loop):
    install root              /usr/local/MiningGuardian
    .env (private secrets)    /usr/local/MiningGuardian/.env  (mode 0600)
    venv                      /usr/local/MiningGuardian/venv
    launcher wrapper          /usr/local/MiningGuardian/bin/intelligence_report_launcher.sh
    logs                      /usr/local/MiningGuardian/logs/
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miningguardian.intelligence-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/MiningGuardian/bin/intelligence_report_launcher.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/usr/local/MiningGuardian</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/usr/local/MiningGuardian/logs/intelligence-report.out.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/MiningGuardian/logs/intelligence-report.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/usr/local/MiningGuardian/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>GUARDIAN_DB</key>
        <string>/usr/local/MiningGuardian/guardian.db</string>
    </dict>
</dict>
</plist>
