# installer/macos-pkg/resources/launchd/com.miningguardian.scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.miningguardian.scanner.plist (PR I-2 / Bucket-3 installer)

  launchd LaunchDaemon for the main Mining Guardian scan loop on macOS.
  Mirrors the systemd unit `mining-guardian.service` running today on the
  ROBS-PC + VPS pair (cutover scope γ — Mini replaces both).

  Install path: copied by postinstall.sh to
      /Library/LaunchDaemons/com.miningguardian.scanner.plist
  and loaded with:
      sudo launchctl bootstrap system <path-to-plist>

  Path + env conventions (matches PR #41 feedback-loop-daemon plist):
    install root              /Library/Application Support/MiningGuardian
    .env (private secrets)    /Library/Application Support/MiningGuardian/.env
    venv                      /Library/Application Support/MiningGuardian/venv
    launcher wrapper          /Library/Application Support/MiningGuardian/bin/scanner_launcher.sh
    logs                      /Library/Application Support/MiningGuardian/logs/

  The wrapper exists because launchd has no EnvironmentFile= equivalent;
  the wrapper sources .env then execs python.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miningguardian.scanner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/Application Support/MiningGuardian/bin/scanner_launcher.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Library/Application Support/MiningGuardian</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/scanner.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/scanner.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Library/Application Support/MiningGuardian/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
