# installer/macos-pkg/resources/launchd/scheduled/com.miningguardian.scheduled.knowledge-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.miningguardian.scheduled.knowledge-backup.plist — D-18 Gap 4 / P-007 (v1.0.3)

  Replaces the cron entry in setup.sh::phase_10_cron for ai/backup_knowledge.py.
  Backs up knowledge.json. CRON_SCHEDULE.md row 4:00 AM.

  Time semantics: StartCalendarInterval fires in the Mac's local timezone.
  Per docs/CRON_SCHEDULE.md "All times are CDT (Mac Mini timezone)" —
  set the Mini's timezone before install.

  No KeepAlive — scheduled jobs run, exit, and wait for the next fire.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miningguardian.scheduled.knowledge-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/Application Support/MiningGuardian/bin/scheduled_job_launcher.sh</string>
        <string>ai/backup_knowledge.py</string>
        <string>knowledge_backup</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Library/Application Support/MiningGuardian</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/scheduled/knowledge_backup.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Application Support/MiningGuardian/logs/scheduled/knowledge_backup.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Library/Application Support/MiningGuardian/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
